# rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
redirect_stderr=true
stopsignal=QUIT
autorestart=true
directory=/root

[program:nginx]
priority=10
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'

[program:web]
priority=10
directory=/usr/local/lib/web/backend
command=/usr/local/lib/web/backend/run.py
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0


[group:x]
programs=xvfb,wm,tint2,pcmanfm,x11vnc,novnc

[program:wm]
priority=15
command=/usr/bin/openbox
environment=DISPLAY=":1",HOME="/root",USER="root"

[program:tint2]
priority=15
directory=%HOME%
command=/usr/bin/tint2
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"

[program:pcmanfm]
priority=15
directory=%HOME%
command=/usr/bin/pcmanfm --desktop --profile LXDE
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"






[program:xvfb]
priority=10
command=/usr/local/bin/xvfb.sh
stopsignal=KILL

[program:x11vnc]
priority=20
command=x11vnc -display :1 -xkb -forever -shared -repeat

[program:novnc]
priority=25
directory=/usr/local/lib/web/frontend/static/novnc
command=bash /usr/local/lib/web/frontend/static/novnc/utils/novnc_proxy --listen 6081
stopasgroup=true

[program:chromium]
priority=30
command=/usr/local/bin/chromium --no-sandbox --disable-gpu --no-first-run --remote-debugging-port=9222 --user-data-dir=/root/.chromium-profile
environment=DISPLAY=":1",HOME="/root"
autorestart=true
startsecs=5

[program:cdp-proxy]
priority=35
command=socat TCP-LISTEN:9223,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:9222
autorestart=true

[program:chrome-devtools-mcp]
priority=35
command=chrome-devtools-mcp
environment=DISPLAY=":1",HOME="/root",PUPPETEER_EXECUTABLE_PATH="/usr/local/bin/chromium"
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
